# flights.service
# Template — install.sh substitutes the placeholders below and writes the
# result to /etc/systemd/system/. Don't copy this file directly.
[Unit]
Description=Flights Overhead Display
Documentation=https://github.com/openskynetwork/opensky-api
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__

# X display (uses the login user's Wayland/X session — see __USER__ above)
Environment=DISPLAY=:0
Environment=XAUTHORITY=__HOME__/.Xauthority

# Uncomment to use fake data (useful for testing the display on the Pi itself)
# Environment=SIMULATE=true

ExecStart=__INSTALL_DIR__/venv/bin/python main.py --fullscreen --interval 45
Restart=always
RestartSec=15
# Allow time for the desktop session to start on boot
ExecStartPre=/bin/sleep 8

[Install]
WantedBy=graphical.target
